# etc/uci-defaults/tch_0060-transformer
#!/bin/sh

if [ ! -f /etc/config/transformer ]; then
	touch /etc/config/transformer
fi

uci show transformer.@main[0] >/dev/null 2>/dev/null
if [ $? -ne 0 ]; then
	uci add transformer main
fi

uci add_list transformer.@main[0].ignore_patterns='^Device%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.SoftwareBank%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.NewSoftwareBank$'
uci add_list transformer.@main[0].ignore_patterns='X_000E50_EnableIPv6$'
uci add_list transformer.@main[0].ignore_patterns='QueueManagement%.Classification%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.ProcessStatus%.'
uci add_list transformer.@main[0].ignore_patterns='%.NumberingPlan%.PrefixInfo%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.ProcessStatus%.Process%.'
uci add_list transformer.@main[0].ignore_patterns='Services%.VoiceService%.{i}%.VoiceProfile%.{i}%.NumberingPlan%.PrefixInfo%.'
uci add_list transformer.@main[0].ignore_patterns='Layer3Forwarding%.'
uci add_list transformer.@main[0].ignore_patterns='QueueManagement%.'
uci add_list transformer.@main[0].ignore_patterns='LANDevice%.{i}%.LANEthernetInterfaceConfig%.'
uci add_list transformer.@main[0].ignore_patterns='LANDevice%.{i}%.WLANConfiguration%.'
uci add_list transformer.@main[0].ignore_patterns='LANDevice%.{i}%.LANHostConfigManagement%.'
uci add_list transformer.@main[0].vendor_patterns='X_TPG_'
uci commit transformer
